# tools/com.tdf2026.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tdf2026.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/joesutheran/tour-de-france-2026/tools/publish.sh</string>
    </array>

    <!-- FALLBACK trigger. The primary trigger is n8n -> com.tdf2026.webhook at 10:00 NZ.
         This fires at 10:00 NZ too (the flip hour) as a belt-and-braces backup in case
         the webhook/n8n path is down. Must match BOUNDARY_HOUR in update_stage.py.
         publish.sh is idempotent — a no-op if nothing changed. A missed run (Mac asleep)
         fires on next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/joesutheran/tour-de-france-2026/logs/tdf-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joesutheran/tour-de-france-2026/logs/tdf-daily.err</string>
</dict>
</plist>
