# tools/com.tdf2026.webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tdf2026.webhook</string>

    <!-- Persistent webhook listener. n8n's daily 10:00 NZ cron POSTs
         {"type":"tdf_daily"} here to trigger the publish. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/joesutheran/tour-de-france-2026/tools/run_listener.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/joesutheran/tour-de-france-2026/logs/webhook.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joesutheran/tour-de-france-2026/logs/webhook.err</string>
</dict>
</plist>
